# systemd/user/update-emacs.service
[Unit]
Description=Pull changes from spacemacs git repository

[Service]
Type=oneshot
ExecStart=/usr/bin/git -C %h/.emacs.d pull
