# systemd/user/update-emacs.timer
[Unit]
Description=Run update-emacs every 30 minutes and on boot

[Timer]
OnBootSec=1min
OnUnitActiveSec=30m

[Install]
WantedBy=timers.target
